# work/ws/devel/etc/catkin/profile.d/10.rosbuild.bat
REM generated from rosbuild/env-hooks/10.rosbuild.bat.em

REM env variables in develspace
set ROS_ROOT=C:/work/ws/src/ros/core/rosbuild
