# work/ws/devel/etc/catkin/profile.d/10.rosconsole.bat
REM generated from rosconsole/env-hooks/10.rosconsole.bat.develspace.em

set ROSCONSOLE_CONFIG_FILE=C:/work/ws/src/ros_comm/tools/rosconsole/config/rosconsole.config
