# work/ws/devel/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "C:/work/ws/src/ros/tools/rosbash/rostcsh"
