# work/ws/devel/etc/catkin/profile.d/15.rosbash.bash
# generated from rosbash/env-hooks/15.rosbash.bash.em

. "C:/work/ws/src/ros/tools/rosbash/rosbash"
