# opt/ros/groovy/x86/etc/catkin/profile.d/15.rosbash.bash
# generated from rosbash/env-hooks/15.rosbash.bash.em

. "$CATKIN_ENV_HOOK_WORKSPACE/share/rosbash/rosbash"
