# opt/ros/groovy/x86/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "$CATKIN_ENV_HOOK_WORKSPACE/share/rosbash/rostcsh"
